# roles/sbot/templates/upstart.conf
description "{{ service }}"

start on started networking
stop on runlevel [016]

limit nofile 1000000 1000000

respawn
respawn limit 5 30

setuid {{ user }}
setgid {{ user }}

chdir /home/{{ user }}

exec sbot server --host {{ inventory_hostname }}
